# crontab
* * * * * python2 /scripts/getevo.py >> /scripts/getevo.log
#let it be
