# crontab
* * * * * python2 /scritps/getevo.py >> /scripts/getevo.log
#let it be
