# crontab
* * * * * python2 /scripts/getevo.py > /proc/1/fd/1 2>/proc/1/fd/2
#let it be
