# crontab
* * * * * /scritps/getevo.php
#let it be
